# supervisord.conf
[supervisord]
nodaemon=true

# Internal API — only Streamlit (same container) talks to it, not exposed publicly.
[program:fastapi]
command=uvicorn app.main:app --host 127.0.0.1 --port 8000
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

# Public UI — served on the single port HF Spaces exposes (7860).
[program:streamlit]
command=streamlit run ui/streamlit_app.py --server.port 7860 --server.address 0.0.0.0 --server.headless true --server.enableCORS false --server.enableXsrfProtection false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
